# decr8.service
[Unit]
Description=start and stop python-telegram-bot decr8_bot 

[Service]
User=decr8
WorkingDirectory=/home/decr8/decr8/decr8
ExecStart=python /home/decr8/decr8/decr8_bot.py

[Install]
WantedBy=multi-user.target
